# deploy/systemd/admin-cron-worker.service
[Unit]
Description=admin-cron worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin-cron
Group=admin-cron
WorkingDirectory=/opt/admin-cron
ExecStart=/opt/admin-cron/bin/admin-cron -f /etc/admin-cron/config.yaml -mode 2
Restart=on-failure
RestartSec=5s
TimeoutStopSec=60s
LimitNOFILE=65535
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
